# .config/systemd/user/mytask.service
[Unit]
Description=My Scheduled Task

[Service]
ExecStart=%h/.config/systemd/user/mytask.sh
