# .config/autostart/cliphist.desktop
[Desktop Entry]
Type=Application
Name=Cliphist
Exec=systemctl --user start cliphist
Terminal=false
