# .config/autostart/xremap.desktop
[Desktop Entry]
Type=Application
Name=Xremap
Exec=systemctl --user start xremap
Terminal=false
